# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command = /usr/local/bin/run
