# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command = /usr/local/bin/run
stdout_logfile = /root/app/logs/supervisor.log
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8
